# tools/turing-smart-screen-python.service
[Unit]
Description=Turing Smart Screen Python HIDdev
After=graphical-session.target network-online.target
Wants=graphical-session.target network-online.target

[Service]
Type=simple

; Recommended usage:
; 1. Copy this file to ~/.config/systemd/user/turing-smart-screen-python.service
; 2. Replace WorkingDirectory and ExecStart with your paths
; 3. Run:
;    systemctl --user daemon-reload
;    systemctl --user enable --now turing-smart-screen-python.service

WorkingDirectory=/absolute/path/to/turing-smart-screen-python-main
ExecStart=/absolute/path/to/python /absolute/path/to/turing-smart-screen-python-main/main.py
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
